# home/endchill/.config/runit/sv/pipewire/run
#!/bin/sh

exec /usr/bin/pipewire &
exec /usr/bin/pipewire-pulse &
exec /usr/bin/wireplumber
